# tee-time-bot.timer
[Unit]
Description=Run Tee Time Booking Bot daily at 7:26 AM (pre-login before 7:30 drop)

[Timer]
OnCalendar=*-*-* 07:26:00
Persistent=true

[Install]
WantedBy=timers.target
